# services/change-wallpaper.timer
[Unit]
Description=Schedule the wallpaper download process

[Timer]
OnActiveSec=3m

[Install]
WantedBy=timers.target
